# nextcloud/32.0/supervisor_conf/default.ini
[supervisord]
nodaemon=true

[program:php-fpm]
command = /usr/sbin/php-fpm84
autostart = true

[program:nginx]
command = /usr/sbin/nginx
autostart = true
stdout_logfile= /dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cron]
command=crond -f
